# rootfs/etc/systemd/system/autostart.service
[Unit]
Description=runs /etc/autostart.sh on system boot

[Service]
ExecStart=sh -c /etc/autostart.sh

[Install]
WantedBy=multi-user.target
